# dist/service/accumulate_energy.service
[Unit]
Description=Program untuk akumulasi energy yg dipakai dan didapat (harvest and enjoy)
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u accumulate_energy.py
WorkingDirectory=/var/lib/sundaya/ehub-bakti/
StandardOutput=inherit
StandardError=inherit
Restart=on-failure
RestartSec=10s
User=root

[Install]
WantedBy=multi-user.target