# dist/service/handle_mosfet.service
[Unit]
Description=Handle Mosfet pengamanan
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u handle_mosfet.py
WorkingDirectory=/var/lib/sundaya/ehub-bakti/
StandardOutput=inherit
StandardError=inherit
User=root

[Install]
WantedBy=multi-user.target