# dist/service/handle_canbus.service
[Unit]
Description=Handle Canbus Communication (PMS and Sensor Arus)
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u handle_canbus.py
WorkingDirectory=/var/lib/sundaya/ehub-bakti/
StandardOutput=inherit
StandardError=inherit
Restart=on-failure
RestartSec=10s
User=root

[Install]
WantedBy=multi-user.target