# dist/service/handle_mosfet.timer
[Unit]
Description=MPPT log data Timer
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/2
Persistent=true
Unit=handle_mosfet.service

[Install]
WantedBy=timers.target