# dist/service/handle_relay.timer
[Unit]
Description=Handle relay timer
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/1
Persistent=true
Unit=handle_relay.service

[Install]
WantedBy=timers.target