# dist/service/keep_alive_dock.timer
[Unit]
Description=MPPT log data Timer
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/1
Persistent=true
Unit=keep_alive_dock.service

[Install]
WantedBy=timers.target