# dist/service/store_log_data.timer
[Unit]
Description=store log data Timer
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=store_log_data.service

[Install]
WantedBy=timers.target