# dist/service/timer_wifi.timer
[Unit]
Description=Program untuk pengecekan jika wifi telah digunakan lebih dari 1 jam
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=timer_wifi.service

[Install]
WantedBy=timers.target